# systemd/vdirsyncer.service
[Unit]
Description=Sync calendars via vdirsyncer (one-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer sync
